# camilladsp-watcher.path
[Path]
Unit=camilladsp-watcher.service
PathChanged=/home/dkdsp/camilladsp-setup

[Install]
WantedBy=multi-user.target
