# services/camilladsp-watcher.path
[Path]
Unit=camilladsp-watcher.service
TriggerLimitBurst=0
PathChanged=/home/dsp-user/dsp/configuration

[Install]
WantedBy=multi-user.target
